# compile/systemd/transpara-ai-civilization-wiki.service
[Unit]
Description=Transpara-AI Civilization Wiki authoring service
Documentation=file:/Transpara/transpara-ai/repos/wiki/compile/REBUILD.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/Transpara/transpara-ai/repos/wiki
Environment=PYTHONUNBUFFERED=1
Environment=CIVWIKI_PYTHON=/Transpara/transpara-ai/repos/wiki/.venv/bin/python3
EnvironmentFile=-%h/.config/transpara-ai-civilization-wiki.env
ExecStart=/bin/sh -lc 'PY="${CIVWIKI_PYTHON:-/Transpara/transpara-ai/repos/wiki/.venv/bin/python3}"; [ -x "$PY" ] || PY=/usr/bin/python3; exec "$PY" /Transpara/transpara-ai/repos/wiki/compile/ingest_server.py 127.0.0.1 8787'
Restart=on-failure
RestartSec=5s
TimeoutStopSec=15s
NoNewPrivileges=true

[Install]
WantedBy=default.target
