# compile/systemd/transpara-ai-civilization-wiki-refresh.timer
[Unit]
Description=Refresh the Transpara-AI Civilization Wiki freshness signal

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
